# playbooks/templates/loadsysctl.service
[Unit]
Description=Run sysctl -p after boot
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 30
ExecStart=/sbin/sysctl -p

[Install]
WantedBy=multi-user.target
